# server/anamnesis-bot.service
[Unit]
Description=Anamnesis Telegram bot (Kobo highlights review + Obsidian capture)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yt2epub
WorkingDirectory=/home/yt2epub/anamnesis
ExecStart=/usr/bin/python3 -u /home/yt2epub/anamnesis/obsidian_bot.py
Restart=on-failure
RestartSec=30
Environment="PYTHONUNBUFFERED=1"
StandardOutput=append:/home/yt2epub/anamnesis/anamnesis_bot.stdout.log
StandardError=append:/home/yt2epub/anamnesis/anamnesis_bot.stderr.log

[Install]
WantedBy=multi-user.target
